# supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=/usr/sbin/apache2ctl -D FOREGROUND

[program:php]
command=php -S 0.0.0.0:8000 -t /var/www/html
